# scanjob/systemd/binder-scanhash.timer
[Unit]
Description=Nightly Binder card-scan hash refresh

[Timer]
# After Scryfall's bulk refresh; job exits immediately when bulk is unchanged.
OnCalendar=*-*-* 03:30:00 UTC
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
